# linux/os/generic/crontab/cleanup-root.sh
#!/bin/bash
echo "########################################################################"
echo "LINUX - OS - DISTRO - DEBIAN-BASED - GENERIC - CRONTAB cleanup root"
echo "########################################################################"


echo "### Crontab cleanup for root"
crontab -u root -r


echo "########################################################################"
echo " "